# local_monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
launchd 設定範本（templates）。安裝步驟見 README「本機監看器」段。

兩個 placeholder：
  __INSTALL_PATH__  → trade_bot/ 的絕對路徑
  __UV_PATH__       → which uv 的結果（通常是 /opt/homebrew/bin/uv）
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradebot.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>__UV_PATH__</string>
        <string>run</string>
        <string>local_monitor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__INSTALL_PATH__</string>

    <!-- 每天 07:00 本機時區跑一次 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- launchd load 時不要立刻跑（避免 install 過程多餘通知）-->
    <key>RunAtLoad</key>
    <false/>

    <!-- 輸出 log（要查紀錄就 tail /tmp/trade_bot_monitor.log） -->
    <key>StandardOutPath</key>
    <string>/tmp/trade_bot_monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/trade_bot_monitor.log</string>
</dict>
</plist>
